# v1-backend/deploy/iot-backend-health.service
[Unit]
Description=IoT Backend Health Check Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/iot-backend-health.sh
